# docker-build/freestyle/fresh.conf
[program:fresh]
autostart=True
autorestart=true
startretries=10000
redirect_stderr=True
command=/usr/local/bin/python3 /app/freestyle/bin/freestyle_config_update.py
user=root
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=20
stdout_logfile=/app/freestyle/logs/fresh.log
directory=/app/freestyle
priority=810